# clearbacon.service
[Unit]
Description=Turn off 7seg on bacon cape
 
[Service]
ExecStart=/usr/bin/clearbacon.sh
Restart=always
 
[Install]
WantedBy=multi-user.target
